# sarah/autostart/gnome-terminal.desktop
[Desktop Entry]
Encoding=UTF-8
Name=gnome-terminal
Comment=gnome-terminal start
Exec=/usr/bin/gnome-terminal
Type=Application

