# debian/system/trackball.service
[Unit]
Description=autostart SlimBlade Trackball

[Service]
Type=simple
ExecStart=/opt/trackball/trackballw.sh

[Install]
WantedBy=multi-user.target
