# debian/autostart/fix4k.desktop
[Desktop Entry]
Type=Application
Name=Fix 4k monitor
Exec=/home/$USER/.xprofile
