# debian/autostart/mydevices.desktop
[Desktop Entry]
Type=Application
Name=My Devices
Exec=/usr/local/sbin/.mydevices
