# debian/autostart/myservices.desktop
[Desktop Entry]
Type=Application
Name=My Services
Exec=/usr/local/sbin/.myservices
